# Quantix-vDC/overlay/etc/inittab
# Quantix-vDC inittab

# Default runlevel
::sysinit:/sbin/openrc sysinit
::sysinit:/sbin/openrc boot
::wait:/sbin/openrc default

# DCUI on tty1 (replaces standard getty)
# The DCUI provides interactive configuration
tty1::respawn:/usr/bin/qx-dcui

# Login terminals on tty2-tty3
tty2::respawn:/sbin/getty 38400 tty2
tty3::respawn:/sbin/getty 38400 tty3

# Serial console (for headless/VM)
ttyS0::respawn:/sbin/getty -L 115200 ttyS0 vt100

# Ctrl-Alt-Delete
::ctrlaltdel:/sbin/reboot

# Shutdown
::shutdown:/sbin/openrc shutdown
