# agent/limiquantix-guest-agent/packaging/debian/prerm
#!/bin/bash
set -e

# Stop the service before removal
if [ "$1" = "remove" ] || [ "$1" = "upgrade" ]; then
    systemctl stop quantix-kvm-agent.service || true
    systemctl disable quantix-kvm-agent.service || true
fi

#DEBHELPER#

exit 0
