# Quantix-vDC/overlay/etc/init.d/nginx
#!/sbin/openrc-run
# =============================================================================
# Quantix-vDC nginx Service
# =============================================================================
# Based on Alpine's nginx init script with additional directory setup.
# =============================================================================

name="nginx"
description="nginx - high performance web server"
extra_commands="checkconfig"
extra_started_commands="reload reopen upgrade"

command="/usr/sbin/nginx"
pidfile="/run/nginx.pid"

depend() {
    need net
    after quantix-controlplane
    use dns logger
}

checkconfig() {
    ebegin "Checking nginx configuration"
    $command -t -q
    eend $?
}

start_pre() {
    # Ensure directories exist
    mkdir -p /var/lib/nginx/logs
    mkdir -p /var/lib/nginx/tmp/client_body
    mkdir -p /var/lib/nginx/tmp/proxy
    mkdir -p /var/lib/nginx/tmp/fastcgi
    mkdir -p /var/log/nginx
    
    # Set permissions (try username, fallback to UID)
    chown -R nginx:nginx /var/lib/nginx /var/log/nginx 2>/dev/null || \
        chown -R 100:101 /var/lib/nginx /var/log/nginx 2>/dev/null || true
    
    # Remove stale PID file if process doesn't exist
    if [ -f "$pidfile" ]; then
        OLD_PID=$(cat "$pidfile" 2>/dev/null)
        if [ -n "$OLD_PID" ] && ! kill -0 "$OLD_PID" 2>/dev/null; then
            ewarn "Removing stale nginx PID file (PID $OLD_PID not running)"
            rm -f "$pidfile"
        fi
    fi
    
    # Validate configuration
    checkconfig || return 1
}

start() {
    ebegin "Starting $name"
    
    # nginx daemonizes itself, so we don't use --background
    start-stop-daemon --start \
        --exec $command \
        --pidfile "$pidfile" \
        -- -g "daemon on;"
    
    eend $?
}

stop() {
    ebegin "Stopping $name"
    
    start-stop-daemon --stop \
        --pidfile "$pidfile" \
        --exec $command \
        --retry QUIT/5/TERM/5/KILL/5
    
    # Clean up PID file
    rm -f "$pidfile"
    
    eend $?
}

reload() {
    checkconfig || return 1
    
    ebegin "Reloading $name configuration"
    
    start-stop-daemon --signal HUP \
        --pidfile "$pidfile"
    
    eend $?
}

reopen() {
    ebegin "Reopening $name log files"
    
    start-stop-daemon --signal USR1 \
        --pidfile "$pidfile"
    
    eend $?
}

upgrade() {
    ebegin "Upgrading $name binary"
    
    start-stop-daemon --signal USR2 \
        --pidfile "$pidfile"
    sleep 3
    
    if [ -f "${pidfile}.oldbin" ]; then
        start-stop-daemon --signal QUIT \
            --pidfile "${pidfile}.oldbin"
    fi
    
    eend $?
}
