# Quantix-vDC/overlay/etc/init.d/quantix-console
#!/sbin/openrc-run
# =============================================================================
# Quantix-vDC Console Service
# =============================================================================
# Launches the interactive DCUI on tty1.
# =============================================================================

name="Quantix-vDC Console"
description="Interactive console DCUI"

depend() {
    need localmount
    after net networking quantix-firstboot quantix-controlplane nginx
    use logger
}

start() {
    ebegin "Starting $name"
    
    # Ensure the DCUI script exists and is executable
    if [ ! -x /usr/bin/qx-dcui ]; then
        ewarn "DCUI not found, showing static info"
        display_static_info
        eend 0
        return
    fi
    
    # Create a wrapper script that respawns DCUI
    cat > /tmp/dcui-wrapper.sh << 'WRAPPER'
#!/bin/sh
while true; do
    /usr/bin/qx-dcui
    # If DCUI exits, wait a moment and restart
    sleep 2
done
WRAPPER
    chmod +x /tmp/dcui-wrapper.sh
    
    # Launch DCUI on tty1
    # Use setsid to create a new session
    setsid sh -c 'exec /tmp/dcui-wrapper.sh < /dev/tty1 > /dev/tty1 2>&1' &
    
    eend 0
}

stop() {
    ebegin "Stopping $name"
    # Kill DCUI and wrapper
    pkill -f "qx-dcui" 2>/dev/null || true
    pkill -f "dcui-wrapper" 2>/dev/null || true
    eend 0
}

display_static_info() {
    # Fallback static display when DCUI is unavailable
    PRIMARY_IP=$(ip -4 addr show scope global 2>/dev/null | grep inet | head -1 | awk '{print $2}' | cut -d/ -f1)
    HOSTNAME=$(hostname 2>/dev/null || cat /etc/hostname)
    
    cat << EOF > /dev/tty1

  ██████╗ ██╗   ██╗ █████╗ ███╗   ██╗████████╗██╗██╗  ██╗
 ██╔═══██╗██║   ██║██╔══██╗████╗  ██║╚══██╔══╝██║╚██╗██╔╝
 ██║   ██║██║   ██║███████║██╔██╗ ██║   ██║   ██║ ╚███╔╝ 
 ██║▄▄ ██║██║   ██║██╔══██║██║╚██╗██║   ██║   ██║ ██╔██╗ 
 ╚██████╔╝╚██████╔╝██║  ██║██║ ╚████║   ██║   ██║██╔╝ ██╗
  ╚══▀▀═╝  ╚═════╝ ╚═╝  ╚═╝╚═╝  ╚═══╝   ╚═╝   ╚═╝╚═╝  ╚═╝

            Quantix-vDC Control Plane Appliance
  ─────────────────────────────────────────────────────────

  Hostname:    ${HOSTNAME}
  IP Address:  ${PRIMARY_IP:-"No network connection"}

  Web Console: https://${PRIMARY_IP:-localhost}/

  ─────────────────────────────────────────────────────────

  Press Ctrl+Alt+F2 for login shell
  
EOF
}
