# quantix-os/overlay/etc/init.d/libvirtd-wrapper
#!/sbin/openrc-run
# ============================================================================
# Quantix-OS libvirtd Wrapper Service
# ============================================================================
# Wraps libvirtd with a startup timeout to prevent boot hangs when KVM
# isn't available (e.g., nested virtualization not enabled).
# ============================================================================

name="libvirtd"
description="Libvirt Virtual Machine Manager (with timeout)"

command="/usr/sbin/libvirtd"
command_args="-d -f /etc/libvirt/libvirtd.conf"
pidfile="/run/libvirtd.pid"

# Key: Don't wait forever for libvirtd to start
start_stop_daemon_args="--background --make-pidfile"

depend() {
    need net
    need dbus
    after firewall virtlogd
    use ceph iscsid nfs rpc.statd
}

start_pre() {
    # Ensure directories exist
    checkpath --directory /run/libvirt
    checkpath --directory /var/log/libvirt
    checkpath --directory /var/lib/libvirt
    
    # Check KVM availability
    if [ ! -c /dev/kvm ]; then
        ewarn "KVM not available - libvirtd will run in emulation mode"
        ewarn "VMs will be slow. Enable nested virtualization for production use."
    fi
    
    # Start virtlogd if not already running
    if ! pgrep -x virtlogd >/dev/null 2>&1; then
        /usr/sbin/virtlogd -d 2>/dev/null || true
    fi
    
    return 0
}

start() {
    ebegin "Starting ${name}"
    
    # Use a timeout to prevent hanging
    # libvirtd may hang if it can't initialize QEMU/KVM properly
    (
        start-stop-daemon --start --quiet \
            --pidfile "${pidfile}" \
            --make-pidfile \
            --background \
            --exec "${command}" -- ${command_args}
    ) &
    
    local start_pid=$!
    local timeout=30
    local count=0
    
    # Wait for libvirtd to start or timeout
    while [ $count -lt $timeout ]; do
        if [ -f "${pidfile}" ]; then
            local pid=$(cat "${pidfile}" 2>/dev/null)
            if [ -n "$pid" ] && kill -0 "$pid" 2>/dev/null; then
                # Check if libvirt is responding
                if virsh list --all >/dev/null 2>&1; then
                    eend 0
                    return 0
                fi
            fi
        fi
        sleep 1
        count=$((count + 1))
    done
    
    # Timeout reached
    ewarn "libvirtd startup timed out (${timeout}s)"
    ewarn "This may indicate missing KVM support"
    
    # Check if it's at least running
    if [ -f "${pidfile}" ] && kill -0 "$(cat ${pidfile})" 2>/dev/null; then
        ewarn "libvirtd is running but may not be fully functional"
        eend 0
        return 0
    fi
    
    eend 1
    return 1
}

stop() {
    ebegin "Stopping ${name}"
    
    if [ -f "${pidfile}" ]; then
        start-stop-daemon --stop --quiet --pidfile "${pidfile}"
        rm -f "${pidfile}"
    fi
    
    eend 0
}

status() {
    if [ -f "${pidfile}" ] && kill -0 "$(cat ${pidfile})" 2>/dev/null; then
        einfo "${name} is running"
        
        if virsh list --all >/dev/null 2>&1; then
            einfo "libvirt API is responsive"
        else
            ewarn "libvirt API is not responding"
        fi
        
        return 0
    fi
    
    einfo "${name} is not running"
    return 1
}
