# quantix-os/overlay/etc/init.d/quantix-console
#!/sbin/openrc-run
# ============================================================================
# Quantix Console Service
# ============================================================================
# Runs the Quantix-OS console GUI on tty1.
# This provides the DCUI (Direct Console User Interface) for local management.
# ============================================================================

name="Quantix Console"
description="Quantix-OS graphical console (DCUI)"

# Run on tty1
command="/usr/bin/qx-console-gui"
command_background="yes"
pidfile="/run/quantix-console.pid"

# Environment for framebuffer rendering
export XDG_RUNTIME_DIR="/run/user/0"
export HOME="/root"

# Slint backend configuration
export SLINT_BACKEND="linuxkms"

depend() {
    need localmount quantix-setup
    after networking
    provide dcui
}

start_pre() {
    # Ensure runtime directory exists
    mkdir -p /run/user/0
    chmod 700 /run/user/0
    
    # Check if the console binary exists
    if [ ! -x /usr/bin/qx-console-gui ]; then
        eerror "Console binary not found at /usr/bin/qx-console-gui"
        return 1
    fi
    
    # Clear the console
    clear > /dev/tty1 2>/dev/null || true
    
    return 0
}

start() {
    ebegin "Starting Quantix Console"
    
    # Run on tty1 with framebuffer access
    start-stop-daemon --start \
        --make-pidfile \
        --pidfile "$pidfile" \
        --background \
        --stdout /var/log/quantix-console.log \
        --stderr /var/log/quantix-console.log \
        --exec "$command" \
        -- < /dev/tty1 > /dev/tty1 2>&1
    
    eend $?
}

stop() {
    ebegin "Stopping Quantix Console"
    
    start-stop-daemon --stop \
        --pidfile "$pidfile" \
        --exec "$command"
    
    eend $?
}

# Fallback to TUI if GUI fails
fallback() {
    ewarn "GUI console failed, falling back to TUI console"
    
    if [ -x /usr/bin/qx-console ]; then
        exec /usr/bin/qx-console < /dev/tty1 > /dev/tty1 2>&1
    else
        ewarn "TUI console not available, dropping to shell"
        exec /bin/sh < /dev/tty1 > /dev/tty1 2>&1
    fi
}
