# quantix-os/overlay/etc/init.d/quantix-setup
#!/sbin/openrc-run
# ============================================================================
# Quantix-OS Setup Service
# ============================================================================
# Creates writable directories needed by services running on read-only root.
# Must run early, before networking and virtualization services.
# ============================================================================

name="Quantix Setup"
description="Create writable directories for Quantix-OS services"

depend() {
    before networking ovsdb-server ovs-vswitchd libvirtd
    after localmount
}

start() {
    ebegin "Setting up Quantix-OS directories"
    
    # Core directories - /quantix stores persistent configuration
    mkdir -p /quantix /data /tmp /run
    chmod 755 /quantix /data
    
    # Console directories (for Slint GUI)
    mkdir -p /run/user/0
    mkdir -p /run/quantix
    chmod 700 /run/user/0
    chmod 755 /run/quantix
    
    # Seatd directories (for KMS session management)
    mkdir -p /run/seatd
    chmod 755 /run/seatd
    
    # OVS directories
    mkdir -p /var/run/openvswitch
    mkdir -p /var/log/openvswitch
    mkdir -p /etc/openvswitch
    chmod 755 /var/run/openvswitch /var/log/openvswitch /etc/openvswitch
    
    # Initialize OVS database if missing
    if [ ! -f /etc/openvswitch/conf.db ]; then
        if [ -f /usr/share/openvswitch/vswitch.ovsschema ]; then
            ovsdb-tool create /etc/openvswitch/conf.db /usr/share/openvswitch/vswitch.ovsschema 2>/dev/null
            einfo "Initialized OVS database"
        fi
    fi
    
    # Libvirt directories
    mkdir -p /var/run/libvirt
    mkdir -p /var/log/libvirt/qemu
    mkdir -p /var/lib/libvirt/images
    mkdir -p /var/lib/libvirt/qemu
    mkdir -p /var/lib/libvirt/dnsmasq
    mkdir -p /var/cache/libvirt/qemu
    mkdir -p /etc/libvirt/qemu/networks
    mkdir -p /etc/libvirt/qemu/autostart
    chmod 755 /var/run/libvirt /var/log/libvirt /var/lib/libvirt
    
    # Check KVM availability
    KVM_AVAILABLE=0
    if [ -c /dev/kvm ]; then
        KVM_AVAILABLE=1
        einfo "KVM hardware acceleration available"
    else
        ewarn "KVM not available - loading kvm modules..."
        # Try to load KVM modules
        modprobe kvm 2>/dev/null || true
        if grep -q "Intel" /proc/cpuinfo 2>/dev/null; then
            modprobe kvm_intel nested=1 2>/dev/null || true
        elif grep -q "AMD" /proc/cpuinfo 2>/dev/null; then
            modprobe kvm_amd nested=1 2>/dev/null || true
        fi
        
        # Check again after loading
        sleep 1
        if [ -c /dev/kvm ]; then
            KVM_AVAILABLE=1
            einfo "KVM modules loaded successfully"
        else
            ewarn "KVM unavailable - nested virtualization may not be enabled"
            ewarn "VMs will use QEMU emulation (slower)"
            ewarn "To enable: Add 'kvm.nested=1' to host's boot params"
        fi
    fi
    
    # Create a flag file for other scripts to check
    if [ "$KVM_AVAILABLE" = "1" ]; then
        touch /run/quantix/kvm_available
    else
        rm -f /run/quantix/kvm_available
    fi
    
    # QEMU directories
    mkdir -p /var/lib/qemu
    mkdir -p /var/cache/qemu
    
    # Quantix directories
    mkdir -p /var/lib/quantix
    mkdir -p /var/log/quantix
    
    # Runtime directories
    mkdir -p /run/lock
    mkdir -p /run/dbus
    mkdir -p /run/udev
    
    # DNS config
    if [ ! -f /run/resolv.conf ]; then
        echo "nameserver 8.8.8.8" > /run/resolv.conf
        echo "nameserver 8.8.4.4" >> /run/resolv.conf
    fi
    
    # Check for first boot
    if [ ! -f /quantix/.setup_complete ]; then
        einfo "First boot detected - installation wizard will run"
    fi
    
    eend 0
}

stop() {
    return 0
}
