# quantix-os/overlay/etc/init.d/wireguard
#!/sbin/openrc-run
# vim: set ft=sh: -*- sh -*-
# Quantix-OS WireGuard Service
# Manages WireGuard VPN interfaces for bastion access

name="WireGuard VPN"
description="WireGuard VPN interfaces"

# Configuration directory
WIREGUARD_DIR="${WIREGUARD_DIR:-/quantix/wireguard}"

depend() {
    need net
    after networking
    before quantix-node
}

start() {
    ebegin "Starting ${name}"
    
    # Check if configuration directory exists
    if [ ! -d "${WIREGUARD_DIR}" ]; then
        ewarn "WireGuard config directory not found: ${WIREGUARD_DIR}"
        ewarn "VPN will not be started (no configuration)"
        eend 0
        return 0
    fi
    
    # Start all configured WireGuard interfaces
    local started=0
    for conf in "${WIREGUARD_DIR}"/*.conf; do
        if [ -f "$conf" ]; then
            local iface=$(basename "$conf" .conf)
            einfo "  Bringing up interface: ${iface}"
            wg-quick up "$conf" && started=$((started + 1))
        fi
    done
    
    if [ $started -eq 0 ]; then
        ewarn "No WireGuard interfaces configured"
    else
        einfo "${started} WireGuard interface(s) started"
    fi
    
    eend 0
}

stop() {
    ebegin "Stopping ${name}"
    
    # Stop all active WireGuard interfaces
    local stopped=0
    for iface in $(wg show interfaces 2>/dev/null); do
        einfo "  Bringing down interface: ${iface}"
        wg-quick down "$iface" 2>/dev/null && stopped=$((stopped + 1))
    done
    
    if [ $stopped -eq 0 ]; then
        einfo "No active WireGuard interfaces"
    else
        einfo "${stopped} WireGuard interface(s) stopped"
    fi
    
    eend 0
}

status() {
    local active=$(wg show interfaces 2>/dev/null)
    if [ -n "$active" ]; then
        einfo "Active WireGuard interfaces:"
        for iface in $active; do
            einfo "  ${iface}:"
            wg show "$iface" 2>/dev/null | sed 's/^/    /'
        done
        return 0
    else
        einfo "No active WireGuard interfaces"
        return 3
    fi
}

reload() {
    ebegin "Reloading ${name}"
    
    # Sync running interfaces with config
    for conf in "${WIREGUARD_DIR}"/*.conf; do
        if [ -f "$conf" ]; then
            local iface=$(basename "$conf" .conf)
            wg syncconf "$iface" <(wg-quick strip "$conf") 2>/dev/null
        fi
    done
    
    eend $?
}
