# quantix-os/overlay/etc/udev/rules.d/99-quantix-graphics.rules
# ============================================================================
# Quantix-OS Graphics Device Rules
# ============================================================================
# Ensures proper permissions for graphics devices used by the console GUI
# ============================================================================

# DRM devices - allow access for Slint LinuxKMS backend
SUBSYSTEM=="drm", KERNEL=="card*", MODE="0666", TAG+="uaccess", TAG+="seat"
SUBSYSTEM=="drm", KERNEL=="renderD*", MODE="0666", TAG+="uaccess"

# Framebuffer devices - allow access for raw framebuffer rendering
SUBSYSTEM=="graphics", KERNEL=="fb*", MODE="0666", TAG+="uaccess", TAG+="seat"

# Input devices - required for Slint evdev input
SUBSYSTEM=="input", KERNEL=="event*", MODE="0666", TAG+="uaccess", TAG+="seat"
SUBSYSTEM=="input", KERNEL=="mouse*", MODE="0666", TAG+="uaccess"

# Virtual console - allow access for TUI fallback
KERNEL=="tty[0-9]*", MODE="0666"
KERNEL=="console", MODE="0600"

# Virtio GPU specific rules
KERNEL=="virtio[0-9]*", SUBSYSTEM=="virtio", ATTR{device}=="0x1050", TAG+="seat"

# Trigger loading of graphics modules on device add
# This helps when hot-plugging or late device enumeration
ACTION=="add", SUBSYSTEM=="pci", ATTR{class}=="0x030000", RUN+="/sbin/modprobe -q virtio_gpu simpledrm"
ACTION=="add", SUBSYSTEM=="pci", ATTR{class}=="0x030200", RUN+="/sbin/modprobe -q virtio_gpu simpledrm"
