# limiquantix-os/overlay/etc/init.d/quantix-node
#!/sbin/openrc-run
# =============================================================================
# Quantix-OS Node Daemon Service
# =============================================================================
# Manages the Quantix node daemon (qx-node).
# =============================================================================

name="Quantix Node Daemon"
description="Quantix-OS hypervisor management daemon"
command="/usr/bin/qx-node"
command_args=""
command_background="yes"
pidfile="/run/quantix-node.pid"
output_log="/var/log/quantix-node.log"
error_log="/var/log/quantix-node.err.log"

depend() {
    need net localmount
    after libvirtd ovsdb-server ovs-vswitchd quantix-firstboot
    provide quantix-node
}

start_pre() {
    # Ensure log directory exists
    mkdir -p /var/log
    
    # Ensure runtime directory exists
    mkdir -p /run/quantix
    mkdir -p /run/user/0
    chmod 700 /run/user/0
    
    # Ensure config directory exists
    mkdir -p /quantix
    
    # Check if binary exists
    if [ ! -x "$command" ]; then
        eerror "Node daemon binary not found: $command"
        return 1
    fi
    
    return 0
}

start() {
    ebegin "Starting $name"
    
    start-stop-daemon --start \
        --background \
        --make-pidfile \
        --pidfile "$pidfile" \
        --stdout "$output_log" \
        --stderr "$error_log" \
        --exec "$command" \
        -- $command_args
    
    eend $?
}

stop() {
    ebegin "Stopping $name"
    
    start-stop-daemon --stop \
        --pidfile "$pidfile" \
        --retry TERM/30/KILL/5
    
    eend $?
}

status() {
    if [ -f "$pidfile" ]; then
        if kill -0 $(cat "$pidfile") 2>/dev/null; then
            einfo "$name is running (PID: $(cat $pidfile))"
            return 0
        fi
    fi
    
    einfo "$name is not running"
    return 3
}
