# limiquantix-os/overlay/etc/init.d/quantix-logrotate
#!/sbin/openrc-run
# =============================================================================
# Quantix-OS Log Rotation Service
# =============================================================================
# Rotates logs to prevent RAM exhaustion (logs are in tmpfs).
# =============================================================================

name="Quantix Log Rotation"
description="Periodic log rotation for Quantix-OS"

depend() {
    need localmount
    after quantix-node
}

start() {
    ebegin "Starting $name"
    
    # Create logrotate config
    mkdir -p /etc/logrotate.d
    
    cat > /etc/logrotate.d/quantix << 'EOF'
/var/log/quantix-node.log {
    size 50M
    rotate 5
    compress
    missingok
    notifempty
    copytruncate
}

/var/log/quantix-node.err.log {
    size 10M
    rotate 10
    compress
    missingok
    notifempty
    copytruncate
}

/var/log/quantix-console.log {
    size 10M
    rotate 3
    compress
    missingok
    notifempty
    copytruncate
}

/var/log/libvirt/*.log {
    size 20M
    rotate 3
    compress
    missingok
    notifempty
    copytruncate
}

/var/log/openvswitch/*.log {
    size 20M
    rotate 3
    compress
    missingok
    notifempty
    copytruncate
}
EOF

    # Run logrotate every 5 minutes via cron-like mechanism
    (
        while true; do
            sleep 300
            logrotate /etc/logrotate.d/quantix 2>/dev/null || true
        done
    ) &
    echo $! > /run/quantix-logrotate.pid
    
    eend 0
}

stop() {
    ebegin "Stopping $name"
    
    if [ -f /run/quantix-logrotate.pid ]; then
        kill $(cat /run/quantix-logrotate.pid) 2>/dev/null || true
        rm -f /run/quantix-logrotate.pid
    fi
    
    eend 0
}
