# quantix-os/overlay/etc/inittab
# /etc/inittab - Quantix-OS init configuration
#
# This inittab is designed for a headless hypervisor with a TUI-based console.
# TTY1 runs the Quantix Console (DCUI - TUI console).
# Other TTYs are disabled for security.

# Default runlevel
::sysinit:/sbin/openrc sysinit
::sysinit:/sbin/openrc boot
::wait:/sbin/openrc default

# TTY1: Quantix Console (DCUI - TUI-based Direct Console User Interface)
# Similar to VMware ESXi's DCUI for initial system configuration
tty1::respawn:/usr/local/bin/qx-console-launcher

# TTY2: Emergency shell (disabled by default, enabled via console menu)
# tty2::respawn:/sbin/getty 38400 tty2

# Serial console for headless servers (disabled by default)
# Enable this on physical hardware with serial console connected:
# ttyS0::respawn:/sbin/getty -L 115200 ttyS0 vt100

# Ctrl-Alt-Del handling
::ctrlaltdel:/sbin/reboot

# Shutdown
::shutdown:/sbin/openrc shutdown
